# etc/supervisord.conf
; PCP Arbitrage — supervisord 配置
;
; 须先启动 supervisord 本进程（否则 supervisorctl 报 sock 不存在）：
;   cd /path/to/pcp_arbitrage && supervisord -c etc/supervisord.conf
;   或：make daemon
;
; 再管理子进程：
;   supervisorctl -c etc/supervisord.conf start pcp_arbitrage
;   supervisorctl -c etc/supervisord.conf stop pcp_arbitrage
;   supervisorctl -c etc/supervisord.conf status
;
; Makefile：make daemon → make up / make down / make status
; 停掉整个 supervisord：make supervisor-shutdown

[unix_http_server]
file = /tmp/pcp_arbitrage_supervisor.sock

[supervisord]
logfile          = data/logs/supervisord.log
logfile_maxbytes = 10MB
logfile_backups  = 3
loglevel         = info
pidfile          = /tmp/pcp_arbitrage_supervisord.pid
nodaemon         = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/pcp_arbitrage_supervisor.sock

[program:pcp_arbitrage
; 使用项目 venv 中的 pcp-arbitrage 入口脚本
command          = %(here)s/../.venv/bin/pcp-arbitrage --config %(here)s/../config.yaml
directory        = %(here)s/..
autostart        = true
autorestart      = true
startsecs        = 5
startretries     = 3
; logging 默认走 stderr：合并到 stdout，只写一份 pcp_arbitrage.log
redirect_stderr     = true
stdout_logfile      = data/logs/pcp_arbitrage.log
stderr_logfile      = data/logs/pcp_arbitrage.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups  = 5
; 优雅停止：先 SIGTERM，等 10s 再 SIGKILL
stopsignal       = TERM
stopwaitsecs     = 10
